# conf/supervisor-celery.ini
[program:celery-worker]

command = celery -A celery_tasks worker -l INFO

directory = /srv/djangoBlog

environment=PATH="/home/ubuntu/.virtualenvs/django-venv/bin"

stdout_logfile=/srv/djangoBlog/log/email.log
stderr_logfile=/srv/djangoBlog/log/email_error.log

autostart=true

autorestart=true
startsecs=0